# dreamshaper.service
[Unit]
Description=LCM Dreamshaper V7 — RK3588 NPU Image Generation Server
Documentation=https://github.com/darkautism/LCM-Dreamshaper-V7-rs
After=network.target

[Service]
Type=simple
ExecStart=%h/.cargo/bin/dreamshaper-cli serve --host 0.0.0.0 --port 8080
Restart=on-failure
RestartSec=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=dreamshaper

[Install]
WantedBy=default.target
